# autostart/ei-web.desktop
[Desktop Entry]
Name=ei-web
Exec=node /home/pi/eiwebHello/index.js
Type=Application
StartupNotify=true
Terminal=true

